# src/terok_sandbox/credentials/proxy/resources/systemd/terok-credential-proxy-socket.service
# terok-credential-proxy-version: {{UNIT_VERSION}}

[Unit]
Description=terok credential proxy (Unix socket mode)

[Service]
Type=simple
ExecStart={{BIN}} --socket-path={{SOCKET_PATH}} --db-path={{DB_PATH}} --routes-file={{ROUTES_PATH}} --ssh-agent-socket-path={{SSH_AGENT_SOCKET_PATH}} --ssh-keys-file={{SSH_KEYS_FILE}}

[Install]
WantedBy=default.target
