# src/terok_sandbox/gate/resources/systemd/terok-gate-socket.service
# terok-gate-version: {{UNIT_VERSION}}
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=terok git gate server (Unix socket mode)

[Service]
Type=simple
ExecStart={{TEROK_GATE_BIN}} --socket-path={{SOCKET_PATH}} --base-path={{GATE_BASE_PATH}} --token-file={{TOKEN_FILE}}

NoNewPrivileges=yes
# Several directives below were absent (AmbientCapabilities=,
# ProtectClock=, ProtectKernelTunables=, ProtectKernelLogs=,
# ProtectKernelModules=, PrivateDevices=) — each implicitly modifies
# CapabilityBoundingSet, which in --user mode on Ubuntu 24.04 hits
# `Failed to drop capabilities: Operation not permitted` (status=218/
# CAPABILITIES) because the user manager lacks CAP_SETPCAP.  Re-add
# behind distro detection in a follow-up.
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
KeyringMode=private
UMask=0077

ProtectControlGroups=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid

PrivateTmp=yes
ProtectSystem=full

# RestrictAddressFamilies=AF_UNIX, PrivateNetwork=yes, and
# IPAddressDeny=any used to live here on the assumption that the
# socket-mode gate speaks only AF_UNIX.  That's wrong: gate
# fetches from upstream forges (github / gitlab / …) over HTTPS,
# so it needs outbound IP.  RestrictAddressFamilies= (seccomp
# BPF) and IPAddressDeny= (cgroup BPF) are enforced even in
# --user mode where namespace-based directives like
# PrivateNetwork= silently downgrade — so the AF_UNIX-only set
# blocked DNS via 127.0.0.53 and the upstream HTTPS connection,
# breaking ``git fetch`` with ``Temporary failure in name
# resolution``.  Match the gate@ inetd variant's allow-list for
# outbound IP (still no AF_INET6 / raw sockets).
RestrictAddressFamilies=AF_UNIX AF_INET
SystemCallFilter=@system-service
SystemCallFilter=~@resources @privileged
MemoryDenyWriteExecute=yes

[Install]
WantedBy=default.target
