# src/terok_sandbox/gate/resources/systemd/terok-gate.socket
# terok-gate-version: {{UNIT_VERSION}}
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=terok git gate socket

[Socket]
ListenStream=127.0.0.1:{{PORT}}
Accept=yes

[Install]
WantedBy=sockets.target
