# src/terok_sandbox/vault/resources/systemd/terok-vault.socket
# terok-vault-version: {{UNIT_VERSION}}

[Unit]
Description=terok vault socket

[Socket]
ListenStream={{SOCKET_PATH}}
SocketMode=0600
ListenStream=127.0.0.1:{{PORT}}

[Install]
WantedBy=sockets.target
