# src/terok_sandbox/vault/resources/systemd/terok-vault-socket.service
# terok-vault-version: {{UNIT_VERSION}}

[Unit]
Description=terok vault (Unix socket mode)

[Service]
Type=simple
ExecStart={{BIN}} "--socket-path={{SOCKET_PATH}}" "--db-path={{DB_PATH}}" "--routes-file={{ROUTES_PATH}}" "--ssh-signer-socket-path={{SSH_SIGNER_SOCKET_PATH}}"

[Install]
WantedBy=default.target
